# place in lib\systemd\system/where-this.service
[Unit]
Description=Program that navigates bluetooth beacons

[Service]
ExecStart=/home/pi/Wedstrijd/init.py
