# place in lib\systemd\system/where-this.timer
[Unit]
Description=Timer to run wherethis two times a day

[Timer]
OnCalendar=11:59
OnCalendar=23:59
Persistent=true
